# parquet_api/deploy/export.cron
# Daily parquet export from Azure SQL, 03:00 UTC.
# Install into the azureuser crontab with: crontab -u azureuser /opt/tfmstats-api/deploy/export.cron
0 3 * * * cd /opt/tfmstats-api && /opt/tfmstats-api/venv/bin/python /opt/tfmstats-api/export_db.py >> /var/log/tfmstats-export.log 2>&1
